# deploy/videox-celery.service
[Unit]
Description=VideoX Celery Worker
After=network.target redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/opt/videox/backend
Environment="PATH=/opt/videox/venv/bin"
Environment="PYTHONPATH=/opt/videox/backend"
ExecStart=/opt/videox/venv/bin/celery -A app.core.celery_app worker \
    --loglevel=info \
    --concurrency=4 \
    --queues=parse,download,default \
    --pool=prefork \
    --logfile=/opt/videox/backend/logs/celery.log \
    --pidfile=/opt/videox/backend/logs/celery.pid
ExecStop=/bin/kill -s TERM $MAINPID
PIDFile=/opt/videox/backend/logs/celery.pid
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
